# quadlet/orca-slicer-api.container
[Unit]
Description=OrcaSlicer API for Bambuddy
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=orca-slicer-api
Image=bambuddy-orca-slicer-api:orca2.3.2
Pull=never
PublishPort=127.0.0.1:3003:3000
Volume=/home/cle/Source/selfhost/orca-slicer-api/data:/app/data
Environment=NODE_ENV=production
Environment=PORT=3000

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
